# supervisord.conf
[supervisord]
nodaemon=true

[program:go_app]
command=/tibet/fml/fml
autorestart=true

[program:uvicorn]
command=uvicorn api:app --host 0.0.0.0 --port 8000
autorestart=true